# upstart/aim.conf
#!upstart
description "aim"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

script
	cd /home/cbsadmin/aimSvr
	exec node aimSvr.js >> /test/logs/aim/aim.log 2>&1
end scrip

pre-start script
	echo "Starting" >> /test/logs/aim/aim.log
end script

pre-stop script
        rm /var/run/aim.pid
        echo "Stopping" >> /test/logs/aim/aim.log
end script
