# upstart/qman.conf
#!upstart
description "qman"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

script
	cd /home/cbsadmin/aimSvr
	exec node aimSvr.js qman >> /test/logs/qman/qman.log 2>&1
end script

pre-start script
	echo "Starting" >> /test/logs/qman/qman.log
end script

pre-stop script
        rm /var/run/qman.pid
        echo "Stopping" >> /test/log/qman/qman.log
end script
